# shutter-buttons.service
[Unit]
Description=Roller Shutter Manual Button Controller
After=network.target shutter-web.service

[Service]
Type=simple
User=akaw
WorkingDirectory=/home/akaw/waveshare-shutter
Environment=SHUTTER_HOME=/home/akaw/waveshare-shutter
ExecStart=/usr/bin/python3 button_daemon.py
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target
